# crontab
* * * * * /usr/local/bin/php /castopod-host/public/index.php scheduled-activities
* * * * * /usr/local/bin/php /castopod-host/public/index.php scheduled-video-clips
