# crontab
* * * * * /usr/local/bin/php /castopod/public/index.php scheduled-activities
