# crontab
* * * * * /usr/local/bin/php /castopod/spark tasks:run >> /dev/null 2>&1
