# docker/production/s6-rc.d/frankenphp/run
#!/command/with-contenv sh
frankenphp run --config /etc/frankenphp/Caddyfile --adapter caddyfile
